# deploy/systemd/worldeval-lab-api.service
[Unit]
Description=WorldEval Controller Lab API
After=network.target

[Service]
Type=simple
User=lissan
Group=lissan
WorkingDirectory=/home/lissan/codex-build-week/WorldEval
Environment=GENESIS_HOST=127.0.0.1
Environment=GENESIS_PORT=18083
Environment=GENESIS_GODOT_EXECUTABLE=/tmp/worldeval-godot-4.5/Godot_v4.5-stable_linux.x86_64
ExecStart=/home/lissan/codex-build-week/WorldEval/.venv/bin/uvicorn genesis_arena.main:app --host 127.0.0.1 --port 18083 --no-access-log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
